# root/etc/services.d/nginx/run
#!/usr/bin/env sh -eu

# Proceed if link exists
docker-nginx-update.sh || [ 100 -eq "$?" ]

exec nginx -g "daemon off;"
