# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh

set -eu

# Proceed if link exists
docker-nginx-update.sh || [ 100 -eq "$?" ]

exec nginx -g "daemon off;"
